# supervisor/app.ini
[program:app]
priority=10
directory=/acme/
command=python /app.py
autorestart=true
stdout_logfile=/var/log/app.log
stderr_logfile=/var/log/app.log