# supervisor/docker_gen.ini
[program:docker_gen]
priority=10
command=/bin/docker-gen -watch -only-exposed -notify 'touch /tmp/force_crt_update' -wait 15s:60s /crt_domains.ini.tmpl /tmp/crt_domains.ini
autorestart=true
stdout_logfile=/var/log/docker_gen.log
stderr_logfile=/var/log/docker_gen.log